# security/config/audit.rules
# Monitor changes to critical system files
-w /etc/passwd -p wa -k identity
-w /etc/shadow -p wa -k identity
-w /etc/group -p wa -k identity
-w /etc/sudoers -p wa -k privilege
-w /etc/ssh/sshd_config -p wa -k ssh
-w /etc/hosts.allow -p wa -k network
-w /etc/hosts.deny -p wa -k network

# Monitor access to executables in /usr/local/bin (common location for backdoors)
-w /usr/local/bin/ -p wa -k local_bin

# Monitor loading/unloading of kernel modules (common rootkit technique)
-w /sbin/insmod -p x -k modules
-w /sbin/rmmod -p x -k modules
-w /sbin/modprobe -p x -k modules
-a always,exit -F arch=b64 -S init_module -S delete_module -k modules

# Monitor system and user management commands
-a always,exit -F arch=b64 -S execve -F path=/usr/bin/sudo -k privilege
-a always,exit -F arch=b64 -S execve -F path=/usr/bin/su -k privilege
-a always,exit -F arch=b64 -S execve -F path=/usr/bin/passwd -k identity
-a always,exit -F arch=b64 -S execve -F path=/usr/bin/useradd -k identity
-a always,exit -F arch=b64 -S execve -F path=/usr/bin/userdel -k identity
-a always,exit -F arch=b64 -S execve -F path=/usr/bin/groupadd -k identity
-a always,exit -F arch=b64 -S execve -F path=/usr/bin/groupdel -k identity

# Monitor network configuration changes
-w /sbin/iptables -p x -k network
-w /sbin/ip6tables -p x -k network
-w /sbin/ebtables -p x -k network
-w /sbin/route -p x -k network
-w /sbin/ifconfig -p x -k network
-w /bin/ip -p x -k network

# Identity & Auth
-w /etc/passwd -p wa -k identity
-w /etc/group -p wa -k identity
-w /etc/shadow -p wa -k identity

# Sudo & SSH
-w /etc/sudoers -p wa -k scope
-w /etc/sudoers.d/ -p wa -k scope
-w /etc/ssh/sshd_config -p wa -k sshd_config

# Systemd & local binaries
-w /etc/systemd/system/ -p wa -k systemd
-w /usr/local/bin/ -p wa -k localbin

# Kernel modules (load/unload)
-a always,exit -F arch=b64 -S init_module -S finit_module -k kernel_module
-a always,exit -F arch=b32 -S init_module -S finit_module -k kernel_module

# Log root execs
-a always,exit -F arch=b64 -S execve -F euid=0 -k exec_as_root
-a always,exit -F arch=b32 -S execve -F euid=0 -k exec_as_root